# adhan.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec ~/my-startups-scripts/adhan/__init__.py